# components/git/etc/profile.d/99git.sh
main ()
{
  git config --global core.editor vi 2> /dev/null
  git config mergetool.prompt false 2> /dev/null
}

main
